# launchd/net.kenliu.pilezero.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  net.kenliu.pilezero.plist — macOS LaunchAgent for the pilezero pipeline.

  This file is a reference template. The recommended installation method is:

      uv run python -m pilezero install-agent [--dry-run]

  That command generates the final plist automatically (correct uv path,
  Homebrew PATH, incoming_dir from ~/.pilezero/config.toml) and loads the
  agent via launchctl. You do not need to edit this file manually.

  If you need to install manually, replace HOMEBREW_BIN with the output of
  `brew --prefix`/bin (e.g. /opt/homebrew/bin or /usr/local/bin), replace
  USERNAME with `whoami`, and replace UUV_PATH with the output of `which uv`.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>net.kenliu.pilezero</string>

  <key>ProgramArguments</key>
  <array>
    <string>UUV_PATH</string>
    <string>run</string>
    <string>python</string>
    <string>-m</string>
    <string>pilezero</string>
    <string>run</string>
    <string>--quiet</string>
  </array>

  <!-- Working directory — replace USERNAME -->
  <key>WorkingDirectory</key>
  <string>/Users/USERNAME/code/pilezero</string>

  <!-- launchd runs with a minimal PATH; add Homebrew so ocrmypdf, tesseract,
       and ghostscript are found. Replace HOMEBREW_BIN with the correct path. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>HOMEBREW_BIN:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- WatchPaths must match incoming_dir in ~/.pilezero/config.toml -->
  <key>WatchPaths</key>
  <array>
    <string>/Users/USERNAME/Dropbox/ScanSnapDestination</string>
  </array>

  <key>StartInterval</key>
  <integer>1800</integer><!-- 30 minutes -->

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/USERNAME/.pilezero/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/.pilezero/launchd.err.log</string>

</dict>
</plist>
